# hatch/striking-clock.service
[Unit]
Description=Striking Clock

[Service]
Type=oneshot
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/.local/bin/striking-clock
